# systemd/buddy-monitor.timer
[Unit]
Description=Run buddy-monitor every 2 minutes

[Timer]
OnUnitActiveSec=2min
AccuracySec=30s
Unit=buddy-monitor.service

[Install]
WantedBy=timers.target
